# templates/default/myusa.upstart.erb
description "Start the myusa unicorn"

start on startup
stop on shutdown

respawn
respawn limit 10 5

script
su <%= @app_user %> -c "/bin/bash <<EOT
export RAILS_ENV=<%= @rails_env %>
export PATH=/opt/rbenv/shims:/opt/rbenv/bin:\"$PATH\"
export SMTP_HOST=<%= @smtp_host %>
export SMTP_PORT=<%= @smtp_port %>
export APP_HOST=<%= @app_host %>
export SMS_NUMBER=<%= @sms_number %>
export ELASTICACHE_ENDPOINT=<%= @elasticache_endpoint %>

<%= @working_dir %>/bin/unicorn -c <%= @working_dir %>/config/unicorn.rb -E staging >> <%= @working_dir %>/log/unicorn.sys.log 2>&1

EOT"
end script
