# templates/default/myusa.upstart.erb
description "Start the myusa unicorn"

start on startup
stop on shutdown

respawn
respawn limit 10 5

setuid <%= @app_user %>

exec . <%= @working_dir %>/config/env.sh && \
<%= @working_dir %>/bin/unicorn -c <%= @working_dir %>/config/unicorn.rb >> <%= @working_dir %>/log/unicorn.sys.log 2>&1
